# crontab
#Add these lines to /etc/crontab
00  12  * * *   pi      php /var/www/html/energie/add_values.php
